# ex00/hw.service
[Unit]
	Description= start hello world script
[Service]
	ExecStart=/usr/bin/hw_script.sh
	Type=simple
	StandardOutput=tty